# instance_configuration/gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
SocketUser=appuser

[Install]
WantedBy=sockets.target